# overlay/etc/init.d/ccapply
#!/sbin/openrc-run

depend() {
    want network-online
}

name="ccapply"
command="/edgi/system/edgi/current/edgi"
command_args="config"